# documents/说明文档/frpc.service
[Unit]
Description=frpc service

[Service]
# ExecStart 设定服务启动的命令
ExecStart=/usr/local/bin/frpc -c /etc/frpc/frpc.ini
